# pi-summoner.service
[Unit]
Description=Pi summoner service

[Service]
ExecStart=/usr/local/bin/pi-summoner
Restart=always
User=pi

[Install]
WantedBy=multi-user.target
